# hosts/iMac-2020/Library/LaunchAgents/local.tailscale-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>local.tailscale-monitor</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/hu/dotfiles/bin/tailscale-healthcheck</string>
	</array>

	<key>StartInterval</key>
	<integer>900</integer>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/hu/.local/state/tailscale-monitor/stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/hu/.local/state/tailscale-monitor/stderr.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>ProcessType</key>
	<string>Background</string>

	<key>Nice</key>
	<integer>1</integer>
</dict>
</plist>
